# frdcsa/cron.d/packager
# Clean up the sandbox directory every now and then
* * * */1 *	jasayne	/var/lib/myfrdcsa/codebases/internal/packager/scripts/clean-up-sandbox.pl
