# frdcsa/cron.d/predator
# Clean up the sandbox directory every now and then
* * * */1 *	jasayne	/var/lib/myfrdcsa/codebases/internal/predator/scripts/clean-up-sandbox.pl
